# doorpi.service.in
[Unit]
Description=DoorPi
Documentation=https://github.com/motom001/DoorPi/wiki
Requires=doorpi.socket
After=doorpi.socket

[Service]
Type=simple
ExecStart=!!prefix!!/bin/doorpi_cli --configfile /etc/!!package!!/!!package!!.ini
Environment=HOME=/var/lib/!!package!!
WorkingDirectory=/var/lib/!!package!!
Restart=always

User=!!package!!
Group=!!package!!
SupplementaryGroups=gpio
DynamicUser=true

CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=true

ProtectSystem=strict
ProtectHome=true
RuntimeDirectory=!!package!!
StateDirectory=!!package!!
LogsDirectory=!!package!!
ConfigurationDirectory=!!package!!

PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

LockPersonality=true
MemoryDenyWriteExecute=false
RestrictRealtime=true
RemoveIPC=true

[Install]
WantedBy=multi-user.target
